# packaging/dlog-main.service
[Unit]
Description=Start the dlog service

[Service]
ExecStart=/usr/bin/dlogutil -v time *:W

[Install]
WantedBy=multi-user.target

